# .config/systemd/user/git-sync-cerebrum.timer
[Unit]
Description=Run git-sync every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
